# .packaging/deb/prerm
#!/bin/bash
set -e
if systemctl is-active --quiet tachibana.service; then
    systemctl stop tachibana.service
fi
systemctl disable tachibana.service || true
systemctl daemon-reload
